# skeleton/etc/inittab
::sysinit:/etc/rcS

::respawn:/sbin/getty 115200 ttyS0

::shutdown:/sbin/ifdown -a
::shutdown:/bin/sync
::shutdown:/bin/umount -a -r

